# vendor/github.com/kyokan/levelsrv/ops/levelsrv.conf
[program:levelsrv]
command=/usr/bin/levelsrv start --db-path /data/chaindata
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/levelsrv.out.log
stderr_logfile=/var/log/supervisor/levelsrv.err.log